# assets/systemd/inir.service
[Unit]
Description=iNiR shell
PartOf=graphical-session.target
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/inir run --session
SuccessExitStatus=143
KillMode=process
KillSignal=SIGTERM
Restart=on-failure
RestartSec=3
StartLimitIntervalSec=30
StartLimitBurst=5
TimeoutStopSec=15

[Install]
WantedBy=graphical-session.target
